# supervisor/brahma.conf
[supervisord]
nodaemon=false
user=root
logfile=/root/.openclaw/workspace/trading-system/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
pidfile=/tmp/brahma_supervisord.pid
childlogdir=/root/.openclaw/workspace/trading-system/logs

[unix_http_server]
file=/tmp/brahma_supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/brahma_supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:ws_guardian]
command=python3 /root/.openclaw/workspace/trading-system/ws_guardian.py
directory=/root/.openclaw/workspace/trading-system
autostart=true
autorestart=true
startretries=999
stopwaitsecs=10
stdout_logfile=/root/.openclaw/workspace/trading-system/logs/ws_guardian.log
stderr_logfile=/root/.openclaw/workspace/trading-system/logs/ws_guardian.log
stdout_logfile_maxbytes=5MB
environment=PYTHONUNBUFFERED="1"

[program:watchdog]
command=bash /root/.openclaw/workspace/trading-system/scripts/watchdog.sh
directory=/root/.openclaw/workspace/trading-system
autostart=true
autorestart=true
startretries=999
stopwaitsecs=10
stdout_logfile=/root/.openclaw/workspace/trading-system/logs/watchdog.log
stderr_logfile=/root/.openclaw/workspace/trading-system/logs/watchdog.log
stdout_logfile_maxbytes=5MB

[program:brahma_daemon]
command=bash /root/.openclaw/workspace/trading-system/scripts/brahma_daemon.sh
directory=/root/.openclaw/workspace/trading-system
autostart=true
autorestart=true
startretries=999
stopwaitsecs=15
stdout_logfile=/root/.openclaw/workspace/trading-system/logs/brahma_daemon.log
stderr_logfile=/root/.openclaw/workspace/trading-system/logs/brahma_daemon.log
stdout_logfile_maxbytes=5MB

[group:brahma]
programs=ws_guardian,watchdog,brahma_daemon
